# ops/macos/com.aipaths.mission-control.local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aipaths.mission-control</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/joaco/Documents/openclaw/repos/aipaths-mission-control-live/node_modules/next/dist/bin/next</string>
    <string>start</string>
    <string>-H</string>
    <string>127.0.0.1</string>
    <string>-p</string>
    <string>3001</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/joaco/Documents/openclaw/repos/aipaths-mission-control-live</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/joaco/Library/Logs/com.aipaths.mission-control.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/joaco/Library/Logs/com.aipaths.mission-control.err.log</string>
</dict>
</plist>
